# src/thelostbottle.socket
[Unit]
Description=The Lost bottle of Rum

[Socket]
ListenStream=5555
Accept=true
NoDelay=true
MaxConnections=200

[Install]
WantedBy=sockets.target
